# deploy/qadbak.service
[Unit]
Description=Qadbak hosting panel (pm2)
Documentation=file:/opt/qadbak/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=QADBAK_DIR=/opt/qadbak
Environment=QADBAK_USER=qadbak
ExecStart=/opt/qadbak/scripts/qadbak-service.sh start
ExecReload=/opt/qadbak/scripts/qadbak-service.sh reload
ExecStop=/opt/qadbak/scripts/qadbak-service.sh stop
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
